# questionnaire.service
[Unit]
Description=Questionnaire API (Node.js)
After=network-online.target
Wants=network-online.target

[Service]
WorkingDirectory=/var/www/okimi
ExecStart=/usr/bin/node /var/www/okimi/server.js
Environment=NODE_ENV=production
EnvironmentFile=-/var/www/okimi/.env
User=ec2-user
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal
# nvm の node を使うなら PATH を追加（例）
# Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/ec2-user/.nvm/versions/node/v20.x/bin

[Install]
WantedBy=multi-user.target
